# opc_azure-vote/app_init.supervisord.conf
[program:sshd]
command=service ssh start
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=root